# .docker/etc/supervisord/services/php-fpm.ini
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autostart=true
autorestart=false
startsecs=0
user=root
redirect_stderr=true
stderr_logfile=/var/log/php-fpm.err.log
stdout_logfile=/var/log/supervisor/php-fpm.log
